# marketing-dashboard-main/ops/systemd/hermes-dashboard.service
[Unit]
Description=Hermes Dashboard (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Template default. Prefer deploying the repo to /opt/hermes-dashboard (or similar)
# instead of user home directories.
WorkingDirectory=/opt/hermes-dashboard

# Secrets belong in an EnvironmentFile, not in the unit file.
# Point this at your deployment-specific env file.
EnvironmentFile=-/etc/hermes-dashboard/hermes-dashboard.env

# Optional: path to a 1Password op:// env template file (non-secret).
Environment=HERMES_OP_ENV_FILE=/etc/hermes-dashboard/hermes-dashboard.op.env
# Optional: 1Password runtime mode (off|auto|required). Default is auto.
Environment=HERMES_1PASSWORD_MODE=auto

ExecStart=/opt/hermes-dashboard/scripts/start-standalone.sh
Restart=on-failure
RestartSec=3

# Run as a dedicated user in template deployments.
User=hermes
Group=hermes

# Hardening baseline. Adjust ReadWritePaths to match your HERMES_DB_PATH/HERMES_STATE_DIR.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hermes-dashboard /var/lib/hermes-dashboard

[Install]
WantedBy=multi-user.target
